# tpl/superdesk-dev/supervisor.conf
{{>superdesk/supervisor.conf}}

[program:client]
command=/bin/sh -c '. {{repo_env}}/bin/activate && exec grunt server --inline'
autostart=true
autorestart=true
stdout_logfile={{logs}}/client.log
redirect_stderr=true
directory={{repo}}/client-core
